# systemd/2602-money-morning.service
[Unit]
Description=2602_money morning briefing
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/2602_money
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/2602_money/.venv/bin/python %h/2602_money/src/jobs/run_morning_briefing.py
StandardOutput=journal
StandardError=journal
