# systemd/2602-money-chatcmd.timer
[Unit]
Description=Run 2602_money chat command poll every minute

[Timer]
OnCalendar=*:0/1
Persistent=true
Unit=2602-money-chatcmd.service

[Install]
WantedBy=timers.target
